# deployment/nexora-learn-backend.service
[Unit]
Description=Nexora Learn Django backend
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/nexora-learn/edutech/edutech_backend
EnvironmentFile=/var/www/nexora-learn/edutech/edutech_backend/.env.production
ExecStart=/var/www/nexora-learn/edutech/edutech_backend/.venv/bin/gunicorn config.wsgi:application \
  --bind 127.0.0.1:8010 \
  --workers 3 \
  --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
