# deployment/nexora-learn-web.service
[Unit]
Description=Nexora Learn Next.js web app
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/nexora-learn/edutech_web
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
